# raspi/backend.service
[Unit]
Description=Go Backend (hackathon_nulabcup)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/workspace/hackathon_nulabcup/backend
ExecStart=/opt/backend/server
EnvironmentFile=%h/workspace/hackathon_nulabcup/backend/.env
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=backend

[Install]
WantedBy=default.target
