# raspi/frontend.service
[Unit]
Description=Next.js Frontend (standalone)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/frontend/app
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

# Static config
Environment=NODE_ENV=production
Environment=PORT=3002
Environment=HOSTNAME=127.0.0.1

# Secrets and runtime config are loaded from .env file
# This file is generated during CD deploy from GitHub Secrets
EnvironmentFile=/opt/frontend/.env

[Install]
WantedBy=default.target
